# mate/root-overlay/etc/rc.local
../../../common/root-overlay/etc/rc.local