# community-qt/root-overlay/etc/skel/.config/autostart/cmst-autostart.desktop
[Desktop Entry]
Categories=Settings;System;Qt;Network;
Comment=QT GUI frontend for connman
Exec=cmst -m
GenericName=Network Configuration
Icon=cmst
Keywords=Network;Wireless;Wi-Fi;Wifi;IP;LAN;Proxy;WAN;Broadband;Bluetooth;vpn;DNS;
Name=Connman UI Setup
Name[de]=Netzwerk-Konfiguration
OnlyShowIn=LXQt;KDE;
StartupNotify=false
Terminal=false
Type=Application
Version=1.0
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel
X-LXQt-Need-Tray=true
