# community-qt/root-overlay/etc/skel/.config/autostart/nm-applet.desktop
[Desktop Entry]
Exec=nm-applet
Name=NetworkManager Applet
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
