# community/root-overlay/etc/rc.local
#!/bin/sh

# Remove possible stale lockfile
rm -f /var/lib/pacman/db.lck &
