# community/root-overlay/etc/rc.local
#!/bin/sh

[ -x /etc/local.d/0-dna_spiral.start ] && /etc/local.d/0-dna_spiral.start &
neofetch >| /etc/issue 

# Remove possible stale lockfile
rm -f /var/lib/pacman/db.lck &
