# community-gtk/live-overlay/etc/rc.local
../../../community/live-overlay/etc/rc.local