# community-qt/root-overlay/etc/rc.local
../../../community/root-overlay/etc/rc.local