# community/live-overlay/etc/local.d/1-remove-openbox-sessions.start
#!/bin/sh

# Openbox is not configured in community ISOs and users can't login as is

rm -f /usr/share/xsessions/openbox*
rm -f $0 &
