# community/live-overlay/etc/local.d/change-machine-id.start
#!/bin/sh

# machine-id can uniquely identify your machine, change it at every boot

dbus-uuidgen  >| /var/lib/dbus/machine-id &
dbus-uuidgen  >| /etc/machine-id &
