# community/live-overlay/etc/local.d/darken-calamares.start
#!/bin/sh

sed -i 's/pkexec //' /usr/share/applications/calamares.desktop &
rm -f $0 &
