# community/live-overlay/etc/local.d/local.stop
#!/bin/bash

[ -x /etc/rc.local.stop ] && /etc/rc.local.stop
