# community-gtk/root-overlay/etc/local.d/local.stop
../../../../community/root-overlay/etc/local.d/local.stop