# community-gtk/root-overlay/etc/local.d/pacman-Syy.start
../../../../community/root-overlay/etc/local.d/pacman-Syy.start