# community-gtk/root-overlay/etc/local.d/theme-root.start
../../../../community/root-overlay/etc/local.d/theme-root.start