# community-qt/root-overlay/etc/local.d/artix-icons.start
../../../../community/root-overlay/etc/local.d/artix-icons.start