# community-qt/root-overlay/etc/local.d/mkinitcpio.start
../../../../community/root-overlay/etc/local.d/mkinitcpio.start