# common/root-overlay/etc/local.d/0-process_local.d.start
#!/bin/sh
# local.d execution in runit/s6

init=$(ls -l /sbin/init | sed 's/.*->\ \(.*\)-init/\1/')

case $init in
    runit|s6)  	for script in /etc/local.d/*.start; do
		    [[ $script =~ local.start ]] && continue
		    [[ $script =~ $0 ]] && continue
		    [[ -x $script ]] && $script
		done
		;;
    *)		;;
esac
