# common/root-overlay/etc/local.d/local.start
#!/bin/bash

[ -x /etc/rc.local ] && /etc/rc.local
