# community/root-overlay/etc/local.d/0-remove-openbox-sessions.start
../../../live-overlay/etc/local.d/0-remove-openbox-sessions.start