# community/root-overlay/etc/local.d/0-remove-openbox-sessions.start
#!/bin/bash

rm -f /usr/share/xsessions/openbox*
rm -f /etc/local.d/0-remove-openbox-sessions.start &
