# community/root-overlay/etc/local.d/1-remove-openbox-sessions.start
../../../live-overlay/etc/local.d/1-remove-openbox-sessions.start