# community/root-overlay/etc/local.d/artix-icons.start
../../../live-overlay/etc/local.d/artix-icons.start