# community/root-overlay/etc/local.d/local.stop
../../../live-overlay/etc/local.d/local.stop