# community/root-overlay/etc/local.d/local.start
../../../live-overlay/etc/local.d/local.start