# community/root-overlay/etc/local.d/mkinitcpio.start
../../../live-overlay/etc/local.d/mkinitcpio.start