# community/root-overlay/etc/local.d/theme-root.start
../../../live-overlay/etc/local.d/theme-root.start